# tools/systemd/mysqlpunk-codex-issue-watch.service
[Unit]
Description=Run Codex only when mySQLPunk GitHub issues changed

[Service]
Type=oneshot
WorkingDirectory=/home/yangchen/projects/mySQLPunk
ExecStart=/home/yangchen/projects/mySQLPunk/tools/codex-issue-watch.sh
